# meta-ov/recipes-apps/sensord/files/sensord.socket
[Unit]
Description=Sensor Daemon for Openvario Sensorboard

[Socket]
ListenStream=%t/sensord.socket
Accept=yes
