# meta-ov/recipes-support/autofs-config/files/boot.automount
[Unit]
Description=Automatically mount the boot partition

[Automount]
Where=/boot
