# meta-ov/recipes-support/autofs-config/files/usb-usbstick.automount
[Unit]
Description=Automatically mount connected USB drives

[Automount]
Where=/usb/usbstick

[Install]
WantedBy=ovmenu-ng.service
