# meta-ov/recipes-support/autofs-config/files/usb-usbstick.mount
[Unit]
Description=Automatically mount connected USB drives

[Mount]
What=/dev/sda1
Where=/usb/usbstick
Options=nodev,noexec,nosuid,noatime,rw,sync
