# recipes-support/autofs-config/files/usb.automount
[Unit]
Description=Automatically mount connected USB drives

[Automount]
Where=/usb

[Install]
WantedBy=ovmenu-ng.service
